# scripts/autostart/ucan-monitor.service.in
[Unit]
Description=Launch UCAN-Monitor at OS start
After=ucan-cansocket.service
[Service]
Type=exec
WorkingDirectory=__project_dir__/bin
ExecStart=__project_dir__/bin/ucan-monitor
Restart=always
RestartSec=2s
TimeoutSec=infinity
[Install]
WantedBy=default.target
