# scripts/autostart/ucan-monitor.desktop
[Desktop Entry]
Name=start-ucan-monitor
Exec=/home/<USER>/start-ucan-monitor.sh
Type=Application
